# deploy/tangra-backup.service
[Unit]
Description=Tangra Backup Server - Centralized Backup Orchestrator
Documentation=https://github.com/go-tangra/go-tangra-backup
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/tangra-backup -c /etc/tangra-backup
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true

# Allow read/write to config, certs, and backup data
ReadWritePaths=/etc/tangra-backup
ReadWritePaths=/var/lib/tangra-backup

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tangra-backup

[Install]
WantedBy=multi-user.target
