# deploy/tangra-backup.service
[Unit]
Description=Tangra Backup Server - Centralized Backup Orchestrator
Documentation=https://github.com/go-tangra/go-tangra-backup
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/tangra-backup -c /etc/tangra-backup
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

# Environment file (create /etc/tangra-backup/env with your values).
# At minimum set: ADMIN_GRPC_ENDPOINT, MODULE_BOOTSTRAP_SECRET, LCM_CA_FINGERPRINT.
EnvironmentFile=-/etc/tangra-backup/env

# ---------------------------------------------------------------------------
# Registration with admin gateway (registration-rework v3.0.0+)
# ---------------------------------------------------------------------------

# Admin gateway gRPC endpoint (plaintext :7787 — registration channel).
Environment=ADMIN_GRPC_ENDPOINT=

# Routable host:port admin uses to dial back into this service.
Environment=GRPC_ADVERTISE_ADDR=0.0.0.0:10100
Environment=HTTP_ADVERTISE_ADDR=0.0.0.0:10101
Environment=FRONTEND_ENTRY_URL=/modules/backup/remoteEntry.js

# Admin's :7787 is plaintext; skip mTLS dial.
Environment=REGISTRATION_INSECURE=1

# ---------------------------------------------------------------------------
# Self-signed cert flow against LCM (registration-rework Phase 1)
# ---------------------------------------------------------------------------

# LCM's bootstrap port. The module dials this with a CSR + secret +
# fingerprint pin and gets back a signed cert + CA bundle.
Environment=LCM_BOOTSTRAP_ENDPOINT=lcm.example.local:9101

# Shared bootstrap secret. MUST match LCM's `module_registration_secret`.
Environment=MODULE_BOOTSTRAP_SECRET=

# Pre-shared SHA-256 fingerprint of LCM's root CA, lowercase hex (64
# chars). Mandatory pin — no TOFU.
Environment=LCM_CA_FINGERPRINT=

# Where cert.Ensure persists the issued ca/server/client cert files.
# Must be writable by the service user.
Environment=CERTS_DIR=/var/lib/tangra-backup/certs

# ---------------------------------------------------------------------------
# Module-specific config
# ---------------------------------------------------------------------------

# Backup storage path
Environment=BACKUP_STORAGE_PATH=/var/lib/tangra-backup

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true

# Allow writes to config dir, cert dir, and backup data dir.
ReadWritePaths=/etc/tangra-backup
ReadWritePaths=/var/lib/tangra-backup

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tangra-backup

[Install]
WantedBy=multi-user.target
